# dockerDeploy/ansible/roles/infra/rabbitmq/files/supervisor-rabbitmq.conf
[program:rabbitmq]
command=/usr/bin/python /opt/rabbitmq/run.py
autostart=true
autorestart=true
redirect_stderr=true
